# haldaemon.upstart
description "HAL Daemon"

start on virtual-filesystems and started SERVICE_syslog=y and started messagebus and started acpid
stop on pld.shutdown-started

respawn
console output

exec /usr/sbin/hald --use-syslog --daemon=no
